# etc/supervisord.d/mysqld.ini
[program:mysqld]
command=/usr/sbin/mysqld --pid-file=/var/run/mysqld/mysqld.pid --user=root
user=root
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
